# supervisord/covid.conf
[program:covid]
autorestart = true
autostart = true
startsecs = 0
directory = /root/covid/
command = /usr/bin/python app.py
startretries = 1
environment=PYTHONPATH="%(here)s"
user=root
